# supervisord.conf
; supervisord.conf — JarvisOS platform
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; Platform API — always starts first
[program:platform-api]
command=python -m platform_api.run
directory=/app
environment=PLATFORM_PORT=8900,PYTHONPATH=/app/src
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ops-detector — Loki poller + Grafana webhook receiver
[program:worker-ops-detector]
command=python -m workers.ops_detector.run
directory=/app
environment=PYTHONPATH=/app/src,LOKI_URL=http://10.10.200.71:3100
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Agent programs are appended by gen_supervisord.py
; Include generated agent configs
[include]
files = /etc/supervisor/conf.d/*.conf
